# app/etc/cron.d/solaranzeige_cron
*     *  *  *  *    /usr/local/bin/php /var/www/solaranzeige/services/regler_auslesen.php
*     *  *  *  *    /usr/local/bin/php /var/www/solaranzeige/services/mqtt.php

20,50 *  *  *  *    /usr/local/bin/php /var/www/solaranzeige/services/wetterdaten.php

55   23  *  *  *    /usr/local/bin/php /var/www/solaranzeige/services/wartung.php
